# build/standalone/supervisord.conf
[supervisord]
logfile = /var/www/html/log/supervisor.log
logfile_maxbytes = 50MB
logfile_backups=0
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = true
minfds = 1024
minprocs = 200
umask = 022
user = root
identifier = supervisor
directory = /tmp
nocleanup = true
childlogdir = /tmp
strip_ansi = false


[program:cron]
command=cron -f
autostart=true
startsecs=10
autorestart=false
stderr_logfile=/var/www/html/log/cron.err.log
stdout_logfile=/var/www/html/log/cron.out.log

[program:apache2]
command=/usr/sbin/apache2ctl -DFOREGROUND -k start
user=www-data
autostart=true
startsecs=10
autorestart=true
killasgroup=true
stopasgroup=true
stderr_logfile=/var/www/html/log/apache2.err.log
stdout_logfile=/var/www/html/log/apache2.out.log

[program:mysql]
command=mysqld_safe
user=www-data
autostart=true
startsecs=10
autorestart=true
killasgroup=true
stopasgroup=true
stderr_logfile=/var/www/html/log/mysql.err.log
stdout_logfile=/var/www/html/log/mysql.out.log

# [program:fail2ban]
# command=/usr/bin/python3 /usr/bin/fail2ban-server -f -s /var/run/fail2ban/fail2ban.sock -p /var/run/fail2ban/fail2ban.pid --logtarget stdout
# startsecs=10

[program:at]
command=atd -f
startsecs=10

# [program:rsyslogd]
# command = rsyslogd -n

# [program:sshd]
# command=/usr/sbin/sshd -D 
