# build/standalone/supervisord.conf
[supervisord]
; supervisor daemon configuration
; http://supervisord.org/configuration.html#supervisord-section-settings
; /proc/1/fd/1 redirect log to the process 1 and File Descriptor 1 (i.e. stdout)
; /proc/1/fd/2 is for stderr
redirect_stderr=true
logfile=/dev/null ; /proc/1/fd/1
logfile_maxbytes = 0 ; 50MB
logfile_backups = 0
loglevel = info
pidfile=/var/run/supervisord.pid
nodaemon = true
minfds = 1024
minprocs = 200
umask = 022
user = root
identifier = supervisor
directory = /tmp
nocleanup = true
childlogdir = /tmp
strip_ansi = false

[unix_http_server]
; required section for apache HTTP server
; http://supervisord.org/configuration.html#unix-http-server-section-settings
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
; add some custom behavior
; http://supervisord.org/configuration.html#rpcinterface-x-section-settings
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
; http://supervisord.org/configuration.html#supervisorctl-section-settings
serverurl = unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

; [include]
; usefull to include more configuration files
; http://supervisord.org/configuration.html#include-section-settings
; files = /etc/supervisor/conf.d/*.conf

; every 'program' sections have similar parameters
; http://supervisord.org/configuration.html#program-x-section-settings
[program:cron]
command=/usr/sbin/cron -f
autostart=false ; start after init.sh
startsecs=10
autorestart=false
; stderr_logfile=/var/log/cron.err.log
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0

[program:apache2]
command=/usr/sbin/apache2ctl -DFOREGROUND -k start
user=www-data
autostart=false ; start after init.sh
startsecs=5
autorestart=true
killasgroup=true
stopasgroup=true
; stderr_logfile=/var/log/apache2.err.log
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0

[program:mysql]
command=mysqld_safe
user=mysql
autostart=true
autorestart=true
startsecs=5
killasgroup=true
stopasgroup=true
; stderr_logfile=/var/log/mysql.err.log
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0

[program:fail2ban]
command = /usr/bin/python3 /usr/bin/fail2ban-server  -fc /etc/fail2ban/
user = root
autostart = false ; start after init.sh
autorestart = true
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0

[program:atd]
command=/usr/sbin/atd -f
autostart = true
autorestart = true
startsecs=2
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0

[program:init]
command = /root/init.sh
user = root
autostart = true
autorestart = false
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0

# [program:rsyslogd]
# command = rsyslogd -n

# [program:sshd]
# command=/usr/sbin/sshd -D 
